# com.pulspm.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pulspm.agent</string>

    <key>ProgramArguments</key>
    <array>
        <!-- install.sh will replace PULSE_BIN_PATH with the real path -->
        <string>PULSE_BIN_PATH</string>
        <string>daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>PULSE_HOME/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>PULSE_HOME/daemon.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
